# docker/s6-overlay-init/services.d/gizmophone/finish
#!/bin/execlineb -P
s6-permafailon 10 2 SIGINT exit