# docker/s6-overlay-init/services.d/gizmophone/run
#!/usr/bin/with-contenv /bin/sh
export HOME=/usr/src/app
exec s6-setuidgid node /usr/src/app/run_main.sh
