# docker/s6-overlay-init/services.d/grove-lcd/run
#!/usr/bin/with-contenv /bin/sh
export HOME=/usr/src/app
exec  /usr/src/app/run_main.sh
