# docker/s6-overlay-init/services.d/cylon-led/run
#!/usr/bin/with-contenv /bin/sh
export HOME=/usr/src/app
cd $HOME
exec node cylon-processes/cylon-led.js
