# docker/s6-overlay-init/services.d/fluidsynth-server/run
#!/usr/bin/with-contenv /bin/sh
# Demo:
cd /data/
fluidsynth -s -a alsa -g 3 /usr/share/sounds/sf2/*.sf2  /data/Dmin\ Chord\ Progression.mid